# .config/autostart/pcmanfm.desktop
[Desktop Entry]
Exec=/usr/bin/pcmanfm --daemon-mode
Name=pcmanfm
Comment=File Manager for Desktop Environment
Type=Application
