# .config/autostart/udiskie.desktop
[Desktop Entry]
Exec=/usr/bin/udiskie -t
Name=udiskie
Type=Application
